# deploy/systemd/backup.timer
[Unit]
Description=Run Car-ID backup daily at 3 AM

[Timer]
OnCalendar=03:00
Unit=backup.service

[Install]
WantedBy=timers.target
